# deploy/ai-agent.service
[Unit]
Description=AI Agent SSH & Troubleshooting Backend
After=network.target

[Service]
User=ai-agent
Group=ai-agent
WorkingDirectory=/opt/ai-agent/backend
EnvironmentFile=/opt/ai-agent/backend/.env
ExecStart=/opt/ai-agent/backend/venv/bin/python main.py
Restart=always
RestartSec=5
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ai-agent/backend/data
LockPersonality=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
